# dotfiles/systemd/develop-autostart.service
[Unit]
Description=Autostart run
After=local-fs.target
After=network.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/sh {{@@ projects.script_dir @@}}/services/runner.sh autostart desktop
Restart=no

[Install]
WantedBy=default.target
