# dotfiles/systemd/nightly-script.timer
[Unit]
Description=Timer of nightly script

[Timer]
OnCalendar=*-*-* 05:30:00
Persistent=false
AccuracySec=5s
Unit=nightly-script.service

[Install]
WantedBy=timers.target
